# root/etc/services.d/nginx/run
#!/bin/sh
exec /usr/local/sbin/nginx -g "daemon off;"
